# scripts/ai.kinglystudio.sunny.env.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  SUNNY launchd user agent: at login, read the Anthropic API key from the
  macOS Keychain (service "sunny-anthropic-api-key") and export it to the
  launchd GUI session via `launchctl setenv ANTHROPIC_API_KEY <value>`.

  After installing this plist, EVERY GUI app the user launches — including
  Sunny.app from Finder or Dock — will see ANTHROPIC_API_KEY in its env.

  Install:
    cp scripts/ai.kinglystudio.sunny.env.plist.template \
       ~/Library/LaunchAgents/ai.kinglystudio.sunny.env.plist
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/ai.kinglystudio.sunny.env.plist

  Remove:
    launchctl bootout gui/$UID/ai.kinglystudio.sunny.env
    rm ~/Library/LaunchAgents/ai.kinglystudio.sunny.env.plist

  Note: this is the OPTIONAL belt-and-braces path. The simpler default is to
  skip this plist entirely and let src-tauri/src/secrets.rs read the Keychain
  directly at Sunny startup.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.kinglystudio.sunny.env</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>KEY="$(/usr/bin/security find-generic-password -s sunny-anthropic-api-key -w 2>/dev/null || true)"; if [ -n "$KEY" ]; then /bin/launchctl setenv ANTHROPIC_API_KEY "$KEY"; fi</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/ai.kinglystudio.sunny.env.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ai.kinglystudio.sunny.env.err.log</string>
</dict>
</plist>
